# deploy/rpi/assets/rud1-agent.service
[Unit]
Description=Rud1 Device Agent
Documentation=https://github.com/arugerdev/rud1-fw
After=network-online.target
Wants=network-online.target
# Restart-rate limiting belongs in [Unit] in systemd >= 230. Anything in
# [Service] is silently ignored by recent systemd (warning in journal).
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
# root is required today because the agent touches WireGuard and USB/IP
# kernel interfaces. When we split those into a privileged helper this can
# drop to a DynamicUser.
User=root
ExecStart=/usr/local/bin/rud1-agent -config /etc/rud1-agent/config.yaml
Restart=always
RestartSec=5

# Observability
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rud1-agent
Environment=RUD1_LOG_LEVEL=info

# Light hardening — kept conservative because the agent needs broad net/fs
# access to manage WireGuard, usbip, system metrics and network snapshots.
NoNewPrivileges=yes
ProtectSystem=full
# The boot identity (registrationCode/PIN) is persisted to the FAT partition
# so it survives an OS reinstall. ProtectSystem=full mounts /boot, /usr and
# /etc read-only for the unit; carve out the writable paths the agent uses:
#   - /boot/firmware (and /boot fallback): the FAT identity file.
#   - /etc/wireguard:                      atomic writes of wg0.conf during
#                                          peer provisioning by the cloud.
#   - /etc/rud1-agent:                     atomic rewrites of config.yaml from
#                                          the setup wizard, NTP/audit endpoints,
#                                          and the cloud→agent desired-config
#                                          applier (cfg.Save() in persist.go).
ReadWritePaths=/boot/firmware /boot /etc/wireguard /etc/rud1-agent
ProtectHome=yes
PrivateTmp=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
# Keep CAP_NET_ADMIN (WireGuard, iptables), CAP_SYS_ADMIN (usbip bind),
# CAP_NET_RAW (ICMP). Everything else is already implicit for root.
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no

# Resource limits — generous but bounded, tuned for a Raspberry Pi 3 / 1GB.
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
